# monitoring_script_timer.timer
[Unit]
Description=Monitor Machine Status

[Timer]
Unit=monitoring_script.service
OnCalendar=*:*:00,30 

[Install]
WantedBy=multi-user.target
