# monitoring_script.service
[Unit]
Description=Run monitoring machine script

[Service]
Type=simple
ExecStart=/bin/monitoring_script.sh

[Install]
WantedBy=multi-user.target
