# monitoring_script_timer.timer
[Unit]
Description=Monitor Machine Status

[Timer]
Unit=minute.timer
Minutes=5 # Every 5 minutes

[Install]
WantedBy=multi-user.target
